# overlay/etc/supervisor/conf.d/sniproxy.conf
[program:sniproxy]
command=/scripts/start_sniproxy.sh
stdout_events_enabled=true
stderr_events_enabled=true
